# crontab
* * * * * . /etc/environment && cd /app/ && vendor/bin/crunz schedule:run >> /var/log/cron.log 2>&1
